# outputs/reports/navigator_layer.json
{
  "name": "Detection Engineering Coverage",
  "versions": {
    "attack": "14",
    "navigator": "4.9.1",
    "layer": "4.5"
  },
  "domain": "enterprise-attack",
  "description": "Auto-generated detection coverage. 2026-04-26T15:53:35 UTC.",
  "filters": {
    "platforms": [
      "Windows",
      "Linux",
      "macOS",
      "Network",
      "PRE",
      "Containers",
      "Office 365",
      "SaaS",
      "Google Workspace",
      "IaaS",
      "Azure AD"
    ]
  },
  "sorting": 0,
  "layout": {
    "layout": "side",
    "aggregateFunction": "average",
    "showID": true,
    "showName": true,
    "showAggregateScores": false,
    "countUnscored": false,
    "expandedSubtechniques": "annotated"
  },
  "hideDisabled": false,
  "techniques": [
    {
      "techniqueID": "T1053.005",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0001: Scheduled Task Creation via schtasks.exe (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1059.001",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0002: Test Sigma Rule (PASS)\nDETECT-2026-0004: Test Sigma Rule (PASS)\nDETECT-2026-0003: Test3 (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    }
  ],
  "gradient": {
    "colors": [
      "#ffffff00",
      "#2EA44F"
    ],
    "minValue": 0,
    "maxValue": 100
  },
  "legendItems": [
    {
      "label": "PASS",
      "color": "#2EA44F"
    },
    {
      "label": "FAIL",
      "color": "#CF222E"
    },
    {
      "label": "Not Verified",
      "color": "#6E7681"
    }
  ],
  "metadata": [],
  "links": [],
  "showTacticRowBackground": true,
  "tacticRowBackground": "#205b8f",
  "selectTechniquesAcrossTactics": true,
  "selectSubtechniquesWithParent": false,
  "selectVisibleTechniques": false
}