# outputs/reports/navigator_layer.json
{
  "name": "Detection Engineering Coverage",
  "versions": {
    "attack": "19",
    "navigator": "4.9.1",
    "layer": "4.5"
  },
  "domain": "enterprise-attack",
  "description": "Auto-generated detection coverage. 2026-05-16T12:08:58 UTC.",
  "filters": {
    "platforms": [
      "Windows",
      "Linux",
      "macOS",
      "Network",
      "PRE",
      "Containers",
      "Office 365",
      "SaaS",
      "Google Workspace",
      "IaaS",
      "Azure AD"
    ]
  },
  "sorting": 0,
  "layout": {
    "layout": "side",
    "aggregateFunction": "average",
    "showID": true,
    "showName": true,
    "showAggregateScores": false,
    "countUnscored": false,
    "expandedSubtechniques": "annotated"
  },
  "hideDisabled": false,
  "techniques": [
    {
      "techniqueID": "T1197",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0005: BITS Job abuse via BITSAdmin or PowerShell (PASS)\nDETECT-2026-0006: BITS Persistence via SetNotifyCmdLine (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1059.001",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0007: PowerShell Encoded Command Execution (PASS)\nDETECT-2026-0008: PowerShell Download Cradle Execution (PASS)\nDETECT-2026-0009: PowerShell Inline Base64 Decode and Execute (PASS)\nDETECT-2026-0010: PowerShell Execution Policy Bypass (PASS)\nDETECT-2026-0011: PowerShell Compressed Payload Decompression and Execute (FAIL)\nDETECT-2026-0012: PowerShell AMSI Bypass (FAIL)\nDETECT-2026-0013: PowerShell Spawned from Suspicious Parent Process (PASS)\nDETECT-2026-0014: PowerShell Stealth Execution Flags Combination (PASS)\nDETECT-2026-0015: PowerShell Reflection-Based Assembly Load (PASS)\nDETECT-2026-0016: PowerShell Script Block or Module Logging Disable (PASS)\nDETECT-2026-0018: PowerShell Known Malicious Cmdlet Invocation (PASS)\nDETECT-2026-0003: Test3 (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1027.010",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0007: PowerShell Encoded Command Execution (PASS)\nDETECT-2026-0009: PowerShell Inline Base64 Decode and Execute (PASS)\nDETECT-2026-0011: PowerShell Compressed Payload Decompression and Execute (FAIL)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1105",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0008: PowerShell Download Cradle Execution (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1685.001",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0012: PowerShell AMSI Bypass (FAIL)\nDETECT-2026-0016: PowerShell Script Block or Module Logging Disable (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1566.001",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0013: PowerShell Spawned from Suspicious Parent Process (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1564.003",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0014: PowerShell Stealth Execution Flags Combination (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1620",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0015: PowerShell Reflection-Based Assembly Load (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1685.005",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0016: PowerShell Script Block or Module Logging Disable (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    },
    {
      "techniqueID": "T1053.005",
      "color": "#2EA44F",
      "comment": "DETECT-2026-0001: Scheduled Task Creation via schtasks.exe (PASS)",
      "enabled": true,
      "score": 100,
      "showSubtechniques": true
    }
  ],
  "gradient": {
    "colors": [
      "#ffffff00",
      "#2EA44F"
    ],
    "minValue": 0,
    "maxValue": 100
  },
  "legendItems": [
    {
      "label": "PASS",
      "color": "#2EA44F"
    },
    {
      "label": "FAIL",
      "color": "#CF222E"
    },
    {
      "label": "Not Verified",
      "color": "#6E7681"
    }
  ],
  "metadata": [],
  "links": [],
  "showTacticRowBackground": true,
  "tacticRowBackground": "#205b8f",
  "selectTechniquesAcrossTactics": true,
  "selectSubtechniquesWithParent": false,
  "selectVisibleTechniques": false
}