# outputs/reports/stats.json
{
  "generated_at": "2026-04-18T08:38:48.592878+00:00",
  "total_rules": 4,
  "total_sigma_rules": 3,
  "total_splunk_rules": 4,
  "total_native_spl_rules": 1,
  "verified_pass": 4,
  "verified_fail": 0,
  "not_verified": 0,
  "pass_rate_pct": 100,
  "pass_rate_color": "brightgreen",
  "by_level": {
    "critical": 1,
    "high": 1,
    "low": 1,
    "medium": 1
  },
  "by_status": {
    "experimental": 1,
    "stable": 1,
    "test": 2
  },
  "by_tactic": {
    "Execution": 4,
    "Persistence": 1
  },
  "rules": [
    {
      "detect_id": "DETECT-2026-0001",
      "title": "Scheduled Task Creation via schtasks.exe",
      "level": "medium",
      "status": "experimental",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "",
      "tactics": [
        "Persistence",
        "Execution"
      ],
      "techniques": [
        "T1053.005"
      ],
      "file_path": "rules/sigma/test1.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0002",
      "title": "Test Sigma Rule",
      "level": "low",
      "status": "test",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "24600422610",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/test2.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0003",
      "title": "Test3",
      "level": "critical",
      "status": "test",
      "source": "native_spl",
      "verdict": "PASS",
      "run_id": "24600422604",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/splunk/test3.spl"
    },
    {
      "detect_id": "DETECT-2026-0004",
      "title": "Test Sigma Rule",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/test4.sigma.yml"
    }
  ]
}