# outputs/reports/stats.json
{
  "generated_at": "2026-05-12T21:07:45.627620+00:00",
  "total_rules": 17,
  "total_sigma_rules": 16,
  "total_splunk_rules": 17,
  "total_native_spl_rules": 1,
  "verified_pass": 16,
  "verified_fail": 1,
  "not_verified": 0,
  "pass_rate_pct": 94,
  "pass_rate_color": "brightgreen",
  "mitre_covered_techniques": 9,
  "mitre_total_techniques": 222,
  "mitre_total_fetched_at": "2026-05-11T15:45:31.975380+00:00",
  "mitre_coverage_pct": 4.1,
  "by_level": {
    "critical": 3,
    "high": 9,
    "low": 1,
    "medium": 4
  },
  "by_status": {
    "experimental": 1,
    "stable": 12,
    "test": 4
  },
  "by_tactic": {
    "Execution": 17,
    "Stealth": 7,
    "Persistence": 3,
    "Defense Impairment": 2,
    "Command & Control": 1,
    "Initial Access": 1
  },
  "rules": [
    {
      "detect_id": "DETECT-2026-0001",
      "title": "Scheduled Task Creation via schtasks.exe",
      "level": "medium",
      "status": "experimental",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "",
      "tactics": [
        "Persistence",
        "Execution"
      ],
      "techniques": [
        "T1053.005"
      ],
      "file_path": "rules/sigma/test1.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0002",
      "title": "Test Sigma Rule",
      "level": "low",
      "status": "test",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "24600422610",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/test2.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0003",
      "title": "Test3",
      "level": "critical",
      "status": "test",
      "source": "native_spl",
      "verdict": "PASS",
      "run_id": "24600422604",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/splunk/test3.spl"
    },
    {
      "detect_id": "DETECT-2026-0004",
      "title": "Test Sigma Rule",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/test4.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0005",
      "title": "BITS Job abuse via BITSAdmin or PowerShell",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25289764033",
      "tactics": [
        "Execution",
        "Persistence",
        "Stealth"
      ],
      "techniques": [
        "T1197"
      ],
      "file_path": "rules/sigma/DETECT-2026-0005_BITS-Job-abuse-via-BITSAdmin-or-PowerShell.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0006",
      "title": "BITS Persistence via SetNotifyCmdLine",
      "level": "critical",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25289764033",
      "tactics": [
        "Execution",
        "Persistence",
        "Stealth"
      ],
      "techniques": [
        "T1197"
      ],
      "file_path": "rules/sigma/DETECT-2026-0006_BITS-Persistence-via-SetNotifyCmdLine.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0007",
      "title": "PowerShell Encoded Command Execution",
      "level": "medium",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25757237398",
      "tactics": [
        "Execution",
        "Stealth"
      ],
      "techniques": [
        "T1059.001",
        "T1027.010"
      ],
      "file_path": "rules/sigma/DETECT-2026-0007_PowerShell-Encoded-Command-Execution.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0008",
      "title": "PowerShell Download Cradle Execution",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25680517513",
      "tactics": [
        "Execution",
        "Command & Control"
      ],
      "techniques": [
        "T1059.001",
        "T1105"
      ],
      "file_path": "rules/sigma/DETECT-2026-0008_PowerShell-Download-Cradle-Execution.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0009",
      "title": "PowerShell Inline Base64 Decode and Execute",
      "level": "critical",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25680517513",
      "tactics": [
        "Execution",
        "Stealth"
      ],
      "techniques": [
        "T1059.001",
        "T1027.010"
      ],
      "file_path": "rules/sigma/DETECT-2026-0009_PowerShell-Base64-Decode-and-Execute.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0010",
      "title": "PowerShell Execution Policy Bypass",
      "level": "medium",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25555176634",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/DETECT-2026-0010_PowerShell-Execution-Policy-Bypass.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0011",
      "title": "PowerShell Compressed Payload Decompression and Execute",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25686499303",
      "tactics": [
        "Execution",
        "Stealth"
      ],
      "techniques": [
        "T1059.001",
        "T1027"
      ],
      "file_path": "rules/sigma/DETECT-2026-0011_PowerShell-Compressed-Payload-Execution.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0012",
      "title": "PowerShell AMSI Bypass",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "FAIL",
      "run_id": "25762181179",
      "tactics": [
        "Execution",
        "Defense Impairment"
      ],
      "techniques": [
        "T1059.001",
        "T1685.001"
      ],
      "file_path": "rules/sigma/DETECT-2026-0012_PowerShell-AMSI-Bypass.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0013",
      "title": "PowerShell Spawned from Suspicious Parent Process",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25556100045",
      "tactics": [
        "Execution",
        "Initial Access"
      ],
      "techniques": [
        "T1059.001",
        "T1566.001"
      ],
      "file_path": "rules/sigma/DETECT-2026-0013_PowerShell-Suspicious-Parent-Process.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0014",
      "title": "PowerShell Stealth Execution Flags Combination",
      "level": "medium",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25757237398",
      "tactics": [
        "Execution",
        "Stealth"
      ],
      "techniques": [
        "T1059.001",
        "T1564.003"
      ],
      "file_path": "rules/sigma/DETECT-2026-0014_PowerShell-Stealth-Execution-Flags.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0015",
      "title": "PowerShell Reflection-Based Assembly Load",
      "level": "high",
      "status": "stable",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25752891371",
      "tactics": [
        "Execution",
        "Stealth"
      ],
      "techniques": [
        "T1059.001",
        "T1620"
      ],
      "file_path": "rules/sigma/DETECT-2026-0015_PowerShell-Reflection-Assembly-Load.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0016",
      "title": "PowerShell Script Block or Module Logging Disable",
      "level": "high",
      "status": "test",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25757983142",
      "tactics": [
        "Execution",
        "Defense Impairment"
      ],
      "techniques": [
        "T1059.001",
        "T1685.001",
        "T1685.005"
      ],
      "file_path": "rules/sigma/DETECT-2026-0016_PowerShell-Logging-Disable.sigma.yml"
    },
    {
      "detect_id": "DETECT-2026-0018",
      "title": "PowerShell Known Malicious Cmdlet Invocation",
      "level": "high",
      "status": "test",
      "source": "sigma",
      "verdict": "PASS",
      "run_id": "25762181179",
      "tactics": [
        "Execution"
      ],
      "techniques": [
        "T1059.001"
      ],
      "file_path": "rules/sigma/DETECT-2026-0018_PowerShell-Invoke-Known-Malicious-Cmdlets.sigma.yml"
    }
  ]
}