# rules/splunk/test1.sigma.spl
META_START
{
  "title": "Scheduled Task Creation via schtasks.exe",
  "detect_id": "DETECT-2026-0001",
  "description": "Detects the creation of a scheduled task via schtasks.exe, a common persistence and lateral movement technique.",
  "sigma_file": "rules/sigma/test1.sigma.yml",
  "level": "medium",
  "status": "test",
  "author": "Bence Marton",
  "date": "2026-04-06",
  "modified": "2026-04-10",

  "convert_time": "2026-04-17T17:18:37+02:00",
  "rule_version": "1.5",
  "git_sha": "4ac3c9f48afd10e0f9626c79255849f3b625b9d3",
  "ci_managed": true,
  "origin": "sigma_to_spl.py",
  "sigma_pipeline": "splunk_sysmon_acceleration",

  "deploy_mode": "alert",
  "index": "sysmon",
  "cron": "*/5 * * * *",
  "earliest": "-5m",
  "latest": "now",
  "severity": "medium",

  "testing enabled": true,
  "runner": "windows-victim",
  "tester": "atomic",
  "atomic tests": [
    {
      "technique": "T1053.005",
      "test_numbers": [
        1,
        2
      ]
    }
  ],

  "references": [
    "https://attack.mitre.org/techniques/T1053/005/"
  ],

  "tags": [
    "attack.persistence",
    "attack.t1053.005",
    "attack.execution"
  ],

  "logsource": {
    "product_category": "os",
    "product": "windows",
    "service": "sysmon",
    "event_type": "process_creation"
  },

  "falsepositives": [
    "Legitimate software installers creating scheduled tasks",
    "Administrative automation scripts"
  ]
}
META_END
---
index=sysmon Image="*\\schtasks.exe" CommandLine="*/create*" | table _time,CommandLine,ParentImage,Image,User,ProcessId
