# rules/sigma/DETECT-2026-0002_win_sysmon_registry-persistence-via-run-keys.sigma.yml
title: Registry Persistence via Run Keys
id: 1e336273-391e-48ab-9d07-2792cfd9bf32
detect_id: DETECT-2026-0002
status: test
description: Detects modifications to registry run keys commonly used for persistence
author: Bence Marton
date: 2026-01-24
modified: 2026-01-26
references:
    - https://example.com
tags:
    - attack.persistence
    - attack.t1547.001

logsource:
    product_category: os
    product: windows
    service: sysmon
    event_type: registry_event

detection:
    # Common persistence registry locations
    selection_registry_paths:
        TargetObject|contains:
            # User run keys
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\'
            # System run keys
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\RunOnce\'
            # Explorer run
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Policies\Explorer\Run\'
            # Winlogon
            - '\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit'
            - '\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell'
            # Active Setup
            - '\SOFTWARE\Microsoft\Active Setup\Installed Components\'

    # Suspicious executable locations in registry values
    selection_suspicious_paths:
        Details|contains:
            - '\AppData\Local\Temp\'
            - '\Users\Public\'
            - '\ProgramData\'
            - '\Windows\Temp\'
            - "%temp%"
            - "%tmp%"
            - "%appdata%"

    # Suspicious execution methods
    selection_suspicious_execution:
        Details|contains:
            - "rundll32"
            - "regsvr32"
            - "powershell"
            - "pwsh"
            - "cmd.exe /c"
            - "cmd /c"
            - "wscript"
            - "cscript"
            - "mshta"
            - "wmic"

    # Additional suspicious patterns
    selection_suspicious_patterns:
        Details|contains:
            - ".vbs"
            - ".vbe"
            - ".js"
            - ".jse"
            - ".wsf"
            - ".wsh"
            - ".ps1"
            - ".bat"
            - ".cmd"
            - "javascript:"
            - ".dll,RunDLL"

    # Filter out common legitimate software
    filter_legitimate_software:
        Details|contains:
            - 'Microsoft\Windows\CurrentVersion\Run\SecurityHealth'
            - "Windows Defender"
            - "OneDrive"
            - "Teams"
            - "Skype"
            - "Chrome"
            - "Firefox"
            - "Adobe"

    filter_legitimate_paths:
        Details|startswith:
            - 'C:\Program Files\'
            - 'C:\Program Files (x86)\'
            - 'C:\Windows\System32\'
            - '"C:\Program Files\'
            - '"C:\Program Files (x86)\'

    # Must be in persistence location AND have suspicious characteristics
    condition: selection_registry_paths and (selection_suspicious_paths or selection_suspicious_execution or selection_suspicious_patterns) and not filter_legitimate_software and not filter_legitimate_paths

fields:
    - EventID
    - TargetObject
    - Details
    - Image
    - User

falsepositives:
    - Legitimate software installations
    - Administrative tools
    - Security software updates

level: high
custom:
    splunk:
        mode: alert
